# build/client-deb/parental-control-client_1.10.1_amd64/DEBIAN/postinst
#!/bin/bash
set -e

# Reload systemd
systemctl daemon-reload

# Enable service
systemctl enable parental-control-client.service

# Restart service if it was already running (upgrade scenario)
if [ "$1" = "configure" ] && [ -n "$2" ]; then
    # This is an upgrade
    echo "Upgrading from version $2 to ${VERSION}..."
    systemctl restart parental-control-client.service || true
else
    # This is a fresh install
    echo ""
    echo "Parental Control Client installed successfully!"
    echo ""
    echo "Next steps:"
    echo "1. Edit /opt/parental-control/appsettings.json"
    echo "2. Set ServerUrl to your server address"
    echo "3. Start service: sudo systemctl start parental-control-client"
    echo "4. Check status: sudo systemctl status parental-control-client"
    echo ""
fi

exit 0
