# build/client-deb/parental-control-client_1.10.1_amd64/DEBIAN/prerm
#!/bin/bash
set -e

# Only stop/disable on removal, not on upgrade
if [ "$1" = "remove" ]; then
    # Stop service if running
    if systemctl is-active --quiet parental-control-client; then
        systemctl stop parental-control-client
    fi

    # Disable service
    systemctl disable parental-control-client.service || true
fi

exit 0
