# build/client-deb/parental-control-client_1.8.0_amd64/DEBIAN/prerm
#!/bin/bash
set -e

# Stop service if running
if systemctl is-active --quiet parental-control-client; then
    systemctl stop parental-control-client
fi

# Disable service
systemctl disable parental-control-client.service || true

exit 0
