# build/tray-deb/parental-control-tray_1.11.1-debug_amd64/etc/xdg/autostart/parental-control-tray.desktop
[Desktop Entry]
Type=Application
Name=Parental Control Tray
Comment=Show remaining time in system tray
Exec=/usr/local/bin/parental-control-tray/ParentalControl.TrayIcon
Terminal=false
Hidden=false
X-GNOME-Autostart-enabled=true
