# build/client-deb/parental-control-client_1.10.1_amd64/etc/systemd/system/parental-control-client.service
[Unit]
Description=Parental Control Client
After=network.target

[Service]
Type=notify
ExecStart=/opt/parental-control/ParentalControl.Client
WorkingDirectory=/opt/parental-control
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=parental-control-client
User=root

[Install]
WantedBy=multi-user.target
